# Milk-V Duo Setup/99-usb-gadget.rules
ACTION=="add", SUBSYSTEM=="usb", ATTR{idVendor}=="1d6b", ATTR{idProduct}=="1009", RUN+="/usr/bin/systemctl start run_usb.service"